# systemd/ai-coach-poll.service
[Unit]
Description=AI Running Coach Garmin Poller
After=network.target

[Service]
Type=oneshot
User=skarambhe
WorkingDirectory=/home/skarambhe/Projects/ai-running-coach
EnvironmentFile=/home/skarambhe/Projects/ai-running-coach/.env
ExecStart=/home/skarambhe/Projects/ai-running-coach/.venv/bin/python -m src.poller
